# scripts/meetmind-bot.service
# systemd unit for the meetmind Slack bot (Linux servers/VMs).
# Install:
#   sudo cp scripts/meetmind-bot.service /etc/systemd/system/
#   # edit WorkingDirectory + User below to match your deployment
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now meetmind-bot
# Logs:  journalctl -u meetmind-bot -f
#
# Tokens are read from the project's .env (loaded by the bot), so no
# EnvironmentFile is needed — just point WorkingDirectory at the project root.

[Unit]
Description=Context Graph Slack bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/E2Eplatform
ExecStart=/usr/bin/python3 -m meetmind.slackbot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
